# overlay/etc/rc.local
mkdir /home/rowan && chown rowan:rowan /home/rowan
su -c "cp -r /etc/skel/.[a-zA-Z0-9]* /home/rowan/" rowan
